# systemd/codex-langfuse-sync.timer
[Unit]
Description=Run Codex-to-Langfuse sync every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=codex-langfuse-sync.service

[Install]
WantedBy=timers.target
