# systemd/codex-langfuse-sync.service
[Unit]
Description=Sync Codex local sessions into Langfuse
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Update EXPORTER_ROOT before installing if you place the project elsewhere.
Environment=EXPORTER_ROOT=/usr/local/langfuse/exporter
WorkingDirectory=/usr/local/langfuse/exporter
Environment=PYTHONUNBUFFERED=1
StateDirectory=codex-langfuse-exporter
ExecStart=/usr/bin/python3 ${EXPORTER_ROOT}/codex_langfuse_sync.py --days 1 --limit 50 --state-file /var/lib/codex-langfuse-exporter/state.json
Nice=10

[Install]
WantedBy=multi-user.target
